# ocarina-os/meta-ocarina/recipes-core/base-files/files/hostname.sh
#!/bin/sh
### BEGIN INIT INFO
# Provides:          hostname
# Required-Start:
# Required-Stop:
# Default-Start:     S
# Default-Stop:
# Short-Description: Set hostname
### END INIT INFO

hostname -F /etc/hostname
